# apps/benchmarks/crontab
# crontab
@hourly ev-benchmarks /root/baseline.json
