# apps/benchmarks/crontab
# crontab
@hourly benchmarks /root/baseline.json
